# bin/.config/systemd/user/power-profile-auto.service
[Unit]
Description=Auto-set power profile based on AC power

[Service]
Type=oneshot
ExecStart=%h/.local/bin/power-profile-auto
